# etc/hotplug.d/iface/99-asterisk
#!/bin/sh
[ "$ACTION" = ifup ] && /etc/init.d/asterisk enabled && /etc/init.d/asterisk restart | sh
