# etc/crontabs/root
*/5 * * * * /root/deviceRestart.sh > /dev/null
