# etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.

ifconfig eth0 txqueuelen 10500
ifconfig wlan0 txqueuelen 7935 mtu 1328
#ifconfig wlan1 txqueuelen 7935 mtu 1328

sleep 25

/sbin/wifi down
/sbin/wifi up

exit 0
